# .config/systemd/user/waybar.service.d/override.conf
[Service]
Type=exec
Slice=app-graphical.slice
