# .config/systemd/user/gammastep.service.d/override.conf
[Service]
Slice=session-graphical.slice

[Service]
TimeoutStopSec=0.2
